# system_resources/system_lord.service
[Unit]
After=network.service

[Service]
ExecStart=/root/Orb-Chat/system_resources/system_lord.sh

[Install]
WantedBy=default.target